# tty.rules
SUBSYSTEM=="tty", KERNEL=="ttyUSB*", ACTION=="add", RUN+="/opt/tty.hello.sh '%E{DEVNAME}'"
